# automount/automount@.service
[Unit]
Description=Mount USB Drive on %i

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/automount.sh mount %i
ExecStop=/usr/local/bin/automount.sh umount %i
